# compat/profile.d/10-runtime-contract.ps1
# Runtime marker used by producer and consumer proofs to confirm that the
# orchestrator-owned profile loader completed its lexical composition.

$env:CFW_PROFILE_COMPOSITION = 'cfw-runtime-v1'
